# deployment/modules/oci/templates/mnt-master.mount
[Mount]
What=${data_volume_device}
Where=${data_volume_mount_point}
